# backend/infra/systemd/studio-frontend.service
[Unit]
Description=Studio Neeklo Frontend SSR
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/studio-neeklo/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1
ExecStart=/usr/bin/node .output/server/index.mjs
Restart=always
RestartSec=5
StandardOutput=append:/var/log/studio-neeklo/frontend.log
StandardError=append:/var/log/studio-neeklo/frontend.log

[Install]
WantedBy=multi-user.target
